# deploy/rpi/beer-tap-controller.service
[Unit]
Description=Beer Tap Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/cybeer/beer-tap-system/rpi-controller
Environment=PYTHONUNBUFFERED=1
Environment=DEVICE_ENV_PATH=/etc/beer-tap/device.env
RuntimeDirectory=beer-tap
RuntimeDirectoryMode=0755
ExecStart=/home/cybeer/.local/share/beer-tap/venvs/controller/bin/python main.py
Restart=always
RestartSec=3
User=cybeer
Group=cybeer

[Install]
WantedBy=multi-user.target
