# deploy/rpi/tap-display-agent.service
[Unit]
Description=Tap Display Agent
After=network-online.target beer-tap-controller.service
Wants=network-online.target beer-tap-controller.service

[Service]
Type=simple
WorkingDirectory=/home/cybeer/beer-tap-system/tap-display-agent
Environment=PYTHONUNBUFFERED=1
Environment=DEVICE_ENV_PATH=/etc/beer-tap/device.env
ExecStartPre=/bin/mkdir -p /run/beer-tap
ExecStart=/home/cybeer/beer-tap-system/.venv-display-agent/bin/uvicorn main:app --host 127.0.0.1 --port 18181
Restart=always
RestartSec=3
User=cybeer
Group=cybeer

[Install]
WantedBy=multi-user.target
